# ksmbd.service.in
[Unit]
Description=ksmbd userspace daemon
Requires=modprobe@ksmbd.service
Wants=network-online.target
After=modprobe@ksmbd.service network.target network-online.target

[Service]
Type=forking
PIDFile=@runstatedir@/ksmbd.lock
ExecStart=@sbindir@/ksmbdctl start
ExecReload=@sbindir@/ksmbdctl reload
ExecStop=@sbindir@/ksmbdctl stop

[Install]
WantedBy=multi-user.target
